# leadville.service
[Unit]
Description=LeadVille Impact Bridge Service
After=network.target

[Service]
Type=simple
User=jrwest
Group=jrwest
WorkingDirectory=/home/jrwest/projects/LeadVille
Environment=PYTHONPATH=/home/jrwest/projects/LeadVille/src
ExecStart=/usr/bin/python3 -m impact_bridge.leadville_bridge
ExecStartPost=/bin/sleep 5
ExecStartPost=/usr/bin/python3 /home/jrwest/projects/LeadVille/scripts/start_fastapi.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target