# leadville.service
[Unit]
Description=LeadVille Impact Bridge Service
After=network.target

[Service]
Type=simple
User=jrwest
Group=jrwest
WorkingDirectory=/home/jrwest/projects/LeadVille
Environment=PYTHONPATH=/home/jrwest/projects/LeadVille/src
ExecStart=/usr/bin/python3 -m impact_bridge.leadville_bridge
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target