# tmp_bt50-prune.service
[Unit]
Description=Bt50 samples prune job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jrwest/projects/LeadVille
ExecStart=/usr/bin/env python3 tools/bt50_prune_wrapper.py --db logs/bt50_samples.db --days 30 --vacuum-weekly --yes
User=jrwest
Group=jrwest

[Install]
WantedBy=multi-user.target
